# config/systemd/user/battery_charge.timer
[Unit]
Description=Notify about low battery charge.

[Timer]
OnUnitActiveSec=10min
OnBootSec=30s

[Install]
WantedBy=timers.target

